# packages/host/deploy/pico-host-update.timer
[Unit]
Description=Periodic pico-host release update check

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
