# packages/host/deploy/pico-host-update.path
[Unit]
Description=Manual pico-host release update trigger

[Path]
PathChanged=/var/lib/pico-host/update-request
PathModified=/var/lib/pico-host/update-request
Unit=pico-host-update.service

[Install]
WantedBy=multi-user.target
