# bridge/deploy/pi-bridge.service
[Unit]
Description=pi-bridge — phone-to-pi WebSocket relay
Documentation=https://github.com/your-org/pi-bridge
After=network-online.target tailscaled.service
Wants=network-online.target
# Restart no more than 5 times in 60s — bursty crash-loops should fail
# loud, not consume the box.
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
User=pi-bridge
Group=pi-bridge
WorkingDirectory=/opt/pi-mobile-workspace/bridge

# Secrets and runtime config live in /etc/pi-bridge/env (chmod 600).
# See deploy/env.example for the keys.
EnvironmentFile=/etc/pi-bridge/env

# Pin HOST so the only public surface is via `tailscale serve`.
# (env file may override these for ops convenience.)
Environment=HOST=127.0.0.1
Environment=PORT=7777
Environment=BRIDGE_DB=/var/lib/pi-bridge/bridge.db
Environment=HOME=/var/lib/pi-bridge
Environment=PI_WORKSPACES_DIR=/var/lib/pi-bridge/workspaces
Environment=NODE_ENV=production

ExecStart=/usr/bin/node --import tsx --no-warnings src/main.ts

Restart=on-failure
RestartSec=2s

# Stop semantics: give the bridge 10s to drain WS connections cleanly
# before SIGKILL. Node handles SIGTERM as a graceful shutdown via the
# hono server's signal trap.
KillMode=mixed
TimeoutStopSec=10s

# ── sandboxing ────────────────────────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
# Writable: the SQLite DB and pi's per-cwd session JSONLs (under $HOME).
ReadWritePaths=/var/lib/pi-bridge
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
# Network: needs IPv4/IPv6 + UNIX (for tailscale). Block raw/packet.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
MemoryDenyWriteExecute=false
# node's V8 JIT requires W^X relaxation; setting true would SIGSEGV.

# ── observability ─────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-bridge

[Install]
WantedBy=multi-user.target
